# crontab
# min   hour    day     month   weekday command
42      *       *       *       *       /app/duckdns.sh
